# etc/supervisor.d/create-user.ini
[program:create-user]
command=/opt/qnib/bin/create-user.sh
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
#stdout_logfile=/var/log/supervisor/%(program_name)s.log
startsecs=0
