# etc/supervisor.d/gatk.ini
[program:gatk]
command=/usr/local/bin/extract_gatk.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
